# ops/systemd/relaynews-origin.service.example
[Unit]
Description=relaynews.ai origin API
After=network.target

[Service]
Type=simple
User=rebase
WorkingDirectory=/home/rebase/apps/relaynews-origin/current
EnvironmentFile=/home/rebase/apps/relaynews-origin/shared/origin.env
ExecStart=/usr/bin/env bash -lc 'cd /home/rebase/apps/relaynews-origin/current && pnpm --filter @relaynews/origin start'
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20
SyslogIdentifier=relaynews-origin

[Install]
WantedBy=multi-user.target
